# deploy/rpg-api.cron
* * * * * www-data cd /var/www/rpg-api/current && /usr/bin/php artisan schedule:run --no-interaction >> /var/log/rpg-api-scheduler.log 2>&1
