# deploy/supervisor-rpg-api.conf
[group:rpg-api]
programs=rpg-api-queue,rpg-api-reverb

[program:rpg-api-queue]
process_name=%(program_name)s
command=/usr/bin/php artisan queue:work redis --queue=rpg-combat --sleep=1 --tries=3 --timeout=60
directory=/var/www/rpg-api/current
user=www-data
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stopwaitsecs=90
environment=APP_ENV="production"
stdout_logfile=/var/log/supervisor/rpg-api-queue.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
redirect_stderr=true

[program:rpg-api-reverb]
process_name=%(program_name)s
command=/usr/bin/php artisan reverb:start --host=127.0.0.1 --port=8081
directory=/var/www/rpg-api/current
user=www-data
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stopwaitsecs=30
environment=APP_ENV="production"
stdout_logfile=/var/log/supervisor/rpg-api-reverb.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
redirect_stderr=true
